# wg-failover.service
[Unit]
Description=WireGuard failover endpoint updater
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wg-failover
