# wg-failover.timer
[Unit]
Description=Run WireGuard failover checks every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
Unit=wg-failover.service

[Install]
WantedBy=timers.target
